# agents/com.rawrbot.planner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rawrbot.planner</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/ben/Projects/work/scripts/run-planner.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<array>
		<!-- 7am Monday -->
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<!-- 7am Friday -->
		<dict>
			<key>Weekday</key>
			<integer>5</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<!-- 7am Saturday -->
		<dict>
			<key>Weekday</key>
			<integer>6</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<!-- 7am Sunday -->
		<dict>
			<key>Weekday</key>
			<integer>0</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>

	<key>StandardOutPath</key>
	<string>/Users/ben/Projects/work/rawr.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/ben/Projects/work/rawr.log</string>

	<key>WorkingDirectory</key>
	<string>/Users/ben/Projects/work</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/ben/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
