# agents/com.rawrbot.auditor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rawrbot.auditor</string>

	<key>ProgramArguments</key>
	<array>
		<string>__REPO_DIR__/scripts/run-auditor.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<array>
		<!-- 6:30am daily (before worker at 7am, so approved plans are ready same day) -->
		<dict>
			<key>Hour</key>
			<integer>6</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</array>

	<key>StandardOutPath</key>
	<string>/dev/null</string>

	<key>StandardErrorPath</key>
	<string>/dev/null</string>

	<key>WorkingDirectory</key>
	<string>__REPO_DIR__</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
